# shopping_scanner/rootfs/etc/services.d/webserver/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Shopping Scanner..."

cd /var/www/html

bashio::log.info "Configuring for Ingress authentication"

# Create config.js
cat > config.js << EOF
window.HA_CONFIG = {
    haUrl: '',
    haToken: ''
};
console.log('Shopping Scanner: Using ingress session authentication');

// Test with a valid endpoint
fetch('/api/config')
    .then(r => {
        console.log('Connection test status:', r.status);
        if (r.ok) {
            console.log('✓ Connected to Home Assistant');
            return r.json();
        }
        throw new Error('HTTP ' + r.status);
    })
    .then(data => console.log('HA Instance:', data.location_name))
    .catch(err => console.error('Connection test failed:', err));
EOF

bashio::log.info "Configuration written"

if [ ! -f "config.js" ]; then
    bashio::log.error "Failed to create config.js"
    exit 1
fi

bashio::log.info "Starting web server on port 8099..."
exec python3 -m http.server 8099