# shopping_scanner/rootfs/etc/services.d/webserver/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Shopping Scanner..."

cd /var/www/html

bashio::log.info "Configuring for Ingress authentication"

# Create config.js - haUrl should be empty since we're using relative paths
cat > config.js << EOF
window.HA_CONFIG = {
    haUrl: '', // Empty - relative URLs work in ingress
    haToken: '' // Empty - ingress uses session cookies
};
console.log('Shopping Scanner: Using ingress session authentication');
console.log('Base URL: (relative paths)');
EOF

bashio::log.info "Configuration written"

if [ ! -f "config.js" ]; then
    bashio::log.error "Failed to create config.js"
    exit 1
fi

bashio::log.info "Starting web server on port 8099..."
exec python3 -m http.server 8099