# shopping_scanner/rootfs/etc/services.d/webserver/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Shopping Scanner with API Proxy..."

if [ -z "${SUPERVISOR_TOKEN}" ]; then
    bashio::log.error "SUPERVISOR_TOKEN is not set!"
    exit 1
fi

bashio::log.info "Supervisor token available"
bashio::log.info "Starting proxy server..."

cd /var/www/html
exec python3 /proxy_server.py