# shopping_scanner/rootfs/etc/services.d/webserver/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Shopping Scanner..."

cd /var/www/html

if [ -z "${SUPERVISOR_TOKEN}" ]; then
    bashio::log.error "SUPERVISOR_TOKEN is not set!"
    exit 1
else
    bashio::log.info "SUPERVISOR_TOKEN is available (length: ${#SUPERVISOR_TOKEN})"
fi

# Create config.js with smart URL detection - inject token directly
cat > config.js << EOF
(function() {
    const isIngress = window.location.pathname.includes('/api/hassio_ingress/');
    const isLocal = window.location.hostname.includes('homeassistant.local') || 
                    window.location.hostname.match(/^192\.168\./) ||
                    window.location.hostname.match(/^10\./) ||
                    window.location.hostname === 'localhost';
    
    let apiUrl;
    if (isIngress || isLocal) {
        apiUrl = '/api';
    } else {
        apiUrl = window.location.origin + '/api';
    }
    
    window.HA_CONFIG = {
        haUrl: apiUrl,
        haToken: '${SUPERVISOR_TOKEN}'
    };
    
    console.log('Shopping Scanner initialized');
    console.log('API URL:', apiUrl);
    console.log('Token length:', window.HA_CONFIG.haToken.length);
    console.log('Context:', { isIngress: isIngress, isLocal: isLocal });
})();
EOF

bashio::log.info "Configuration written to config.js with smart URL detection"

if [ -f "config.js" ]; then
    bashio::log.info "config.js created successfully"
else
    bashio::log.error "Failed to create config.js"
    exit 1
fi

bashio::log.info "Starting web server on port 8099..."
exec python3 -m http.server 8099