# shopping_scanner/rootfs/etc/services.d/webserver/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Shopping Scanner..."

cd /var/www/html

if [ -z "${SUPERVISOR_TOKEN}" ]; then
    bashio::log.error "SUPERVISOR_TOKEN is not set!"
else
    bashio::log.info "SUPERVISOR_TOKEN is available (length: ${#SUPERVISOR_TOKEN})"
fi

# Create config.js - try using the Ingress API endpoint
cat > config.js << EOF
window.HA_CONFIG = {
    haUrl: '../api/hassio_ingress',
    haToken: '${SUPERVISOR_TOKEN}'
};
EOF

bashio::log.info "Configuration written to config.js"
bashio::log.info "Using Ingress API endpoint"

if [ -f "config.js" ]; then
    bashio::log.info "config.js created successfully"
else
    bashio::log.error "Failed to create config.js"
fi

bashio::log.info "Starting web server on port 8099..."
exec python3 -m http.server 8099